# fix-time-sync.service
[Unit]
Description=Fix system clock drift after snapshots
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/fix-time-sync-final.sh

[Install]
WantedBy=multi-user.target
